# etc/deepseek-balance.service
[Unit]
Description=DeepSeek Balance Tracker — poll and record balance snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/deepseek-balance/secrets.env
# PATH: adjust to your distribution if the defaults below don't cover all
# dependencies (curl, jq, sqlite3 must be on PATH).
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# ExecStart: adjust this path to where you cloned the repository.
ExecStart=%h/code/deepseek-balance-tracker/bin/deepseek-balance-poll

StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepseek-balance
