# etc/deepseek-balance.service
[Unit]
Description=DeepSeek Balance Tracker — poll and record balance snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/deepseek-balance/secrets.env
ExecStart=/home/kelvin/code/deepseek-balance-tracker/bin/poll-balance

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepseek-balance
