# etc/deepseek-balance.timer
[Unit]
Description=DeepSeek Balance Tracker — hourly poll

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
