# etc/deepseek-balance.service
[Unit]
Description=DeepSeek Balance Tracker — poll and record balance snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/deepseek-balance/secrets.env
Environment=PATH=/run/current-system/sw/bin:%h/.nix-profile/bin:/etc/profiles/per-user/%U/bin
ExecStart=%h/code/deepseek-balance-tracker/bin/poll-balance

StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepseek-balance
